# root/clash/files/clash.service
[Unit]
Description=Clash
After=network.target network-online.target
Requires=network-online.target
Wants=network-online.target
[Service]
Type=simple
ExecStartPre=/bin/sleep 2
ExecStart=/opt/clash/bin/start_clash_debian.sh
[Install]
WantedBy=multi-user.target